# systemd/.config/systemd/user/Notes.service
[Unit]
Descripton=Sync notes

[Service]
Type=oneshot
ExecStart="/home/fox/Notes/commit.sh"