# systemd/.config/systemd/user/task.timer
[Unit]
Description=Sync taskwarrior

[Timer]
Persistent=true
OnBootSec=30s
OnCalendar=hourly

[Install]
WantedBy=timers.target