# systemd/.config/systemd/user/task.service
[Unit]
Descripton=Sync task

[Service]
Type=oneshot
ExecStart="/usr/bin/task sync"