# systemd/.config/systemd/user/urlwatch-community.timer
[Unit]
Description=Check urlwatch

[Timer]
Persistent=true
OnBootSec=30s
OnCalendar=hourly

[Install]
WantedBy=timers.target