# systemd/.config/systemd/user/xorg.target
[Unit]
Description=The basic Xorg Server
