# util/etc/udev/rules.d/99-hdmi.rules
KERNEL=="card0", ACTION=="change", RUN+="/home/fox/bin/hdmi-plugged"
